# package/gnome2/desktop-file-utils/update-desktop-database.cron
#!/bin/sh

PREFIX="D_prefix"
echo "Updating desktop database(s)"
/$PREFIX/bin/update-desktop-database -q /{opt/*,usr}/share/applications
